# install-and-deploy/roles/oracle-jdk8/templates/profile.d/java.sh
#!/bin/bash

export LANG=ja_JP.UTF-8
export JAVA_HOME={{ local_home }}/{{jdk_extract_dir}}
export PATH=$PATH:$JAVA_HOME/bin
